# php-fpm/crontab
* * * * * echo "Save mySQL" >> /application/var/log/cron.l 2>&1
* * * * * /backup.sh 2>&1
