# php-fpm/crontab
* * * * * echo "Save mySQL" >> /application/var/log/cron.l 2>&1
* * * * * /backup.sh 2>&1
0 3,7,11,15,19,23 * * * /application/bin/console dump-bdd > /etl.log
