# php-fpm/crontab
* * * * * echo "Save mySQL" >> /var/log/cron.log 2>&1
* * * * * /backup.sh 2>&1
